# docker/etc/profile.d/enable-thirdparty.sh
if [ -d /depends/thirdparty/bin ]; then
    export PATH=/depends/thirdparty/bin:$PATH
fi

if [ -d /opt/maven/bin ]; then
    export PATH=/opt/maven/bin:$PATH
fi
