# docker/etc/profile.d/enable-rh.sh
if [ -d '/opt/rh/devtoolset-8' ] ; then
    source /opt/rh/devtoolset-8/enable
fi

if [ -d '/opt/rh/rh-python38' ]; then
    source /opt/rh/rh-python38/enable
fi

if [ -d '/opt/rh/rh-git218' ] ; then
    source /opt/rh/rh-git218/enable
fi
