# docker/etc/profile.d/enable-thirdparty.sh
if [ -d /depends/thirdparty/bin ]; then
    export PATH=/depends/thirdparty/bin:$PATH
fi
